# etc/supervisor/conf.d/elasticsearch.conf
[program:elasticsearch]
command=/opt/logstash_deploy/apps/elasticsearch-0.90.10/bin/elasticsearch -f -p /opt/logstash_deploy/log/elasticsearch/elasticsearch.pid
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
autostart=true                ; start at supervisord start (default: true)
;redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=/opt/logstash_deploy/log/elasticsearch/stdout.log
stdout_logfile_maxbytes=1MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10     ; # of stdout logfile backups (default 10)
stderr_logfile=/opt/logstash_deploy/log/elasticsearch/stderr.log
stderr_logfile_maxbytes=1MB   ; max # logfile bytes b4 rotation (default 50MB)
stderr_logfile_backups=10     ; # of stderr logfile backups (default 10)
